# crontab
# Caps Edge data refresh schedule
# Runs at 1:00 PM, 7:00 PM, and 11:00 PM ET

SHELL=/bin/bash
PATH=/usr/local/bin:/usr/bin:/bin
PYTHONPATH=/app
DATA_DIR=/app/data

# 1:00 PM ET (13:00)
0 13 * * * root cd /app && python scripts/refresh.py >> /var/log/caps-edge-refresh.log 2>&1

# 7:00 PM ET (19:00)
0 19 * * * root cd /app && python scripts/refresh.py >> /var/log/caps-edge-refresh.log 2>&1

# 11:00 PM ET (23:00)
0 23 * * * root cd /app && python scripts/refresh.py >> /var/log/caps-edge-refresh.log 2>&1

# Empty line required at end of crontab
